# contrib/vector_service/systemd/debgpt-vector.service
[Unit]
Description=DebGPT vector service (FastAPI)
After=network.target docker.service
Requires=docker.service

[Service]
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/debgpt-vector
Environment="PATH=/home/YOUR_USER/debgpt-vector/.venv/bin"
ExecStart=/home/YOUR_USER/debgpt-vector/.venv/bin/uvicorn debgpt.vector_service.app:app --host 127.0.0.1 --port 8000 --workers 2
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
